# udev/99-safecast.rules
ACTION=="add", KERNEL=="ttyUSB0", RUN+="/bin/bash -c '/usr/sbin/safecast-autoprog'"